# templates/systemd/system/telegraf.service.d/override.conf
[Service]
User={{ telegraf_runas_user }}
